# cookbooks/kubernetes/files/default/kubelet
#!/sbin/openrc-run

command="/usr/bin/kubelet"
pidfile="/run/${RC_SVCNAME}.pid"
command_args=" ${DAEMON_ARGS}"
LOGFILE="${LOGFILE:-/var/log/${RC_SVCNAME}.log}"
start_stop_daemon_args="--background \
    --stderr \"${LOGFILE}\" --stdout \"${LOGFILE}\" --make-pidfile --pidfile ${pidfile}"


