# cookbooks/kubernetes/templates/default/kube-proxy.service.erb
[Unit]
Description=kube-proxy
Documentation=man page

[Service]
EnvironmentFile=-/etc/default/%p
ExecStart=/usr/bin/kube-proxy $DAEMON_ARGS
Restart=always
RestartSec=10s

[Install]
WantedBy=multi-user.target
