# cookbooks/kubernetes/templates/default/kubelet.service.erb
[Unit]
Description=kubelet
Documentation=man page

[Service]
EnvironmentFile=/root/.aws/env
EnvironmentFile=-/etc/default/%p
ExecStart=/usr/bin/kubelet $DAEMON_ARGS
Restart=always
RestartSec=10s

[Install]
WantedBy=multi-user.target
