# cookbooks/kubernetes/templates/default/kube-apiserver.service.erb
[Unit]
Description=kube-apiserver
Documentation=man page

[Service]
EnvironmentFile=/root/.aws/env
EnvironmentFile=/etc/default/%p
ExecStart=/usr/bin/kube-apiserver $DAEMON_ARGS
Restart=always
RestartSec=10s

[Install]
WantedBy=multi-user.target
